# batt_checker.timer
[Unit]
Description=Check battery

[Timer]
OnBootSec=300
OnUnitActiveSec=900

[Install]
WantedBy=timers.target
